# node_ssh/nst.cron
NTS_PATH = "/var/planetlab/test/nts"

*/2 * * * * cd $NTS_PATH; python nts.py;



